# config/systemd/user/drivebackup.service
[Unit]
Description=Backup to google drive

[Service]
Type=simple
ExecStart=%h/bin/drivebackup.sh -a

[Install]
WantedBy=basic.target

