# config/systemd/user/emptytrash.service
[Unit]
Description=Empty the trash

[Service]
Type=oneshot
ExecStart=trash-empty 30
