# config/systemd/user/emptytrash.timer
[Unit]
Description=Empty the (30 day old) trash once a day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
